# _F9K1116_WW_1.03.08.bin.extracted/fmk/rootfs/etc/inittab
::sysinit:/etc/rcS
console::respawn:/bin/sh
